# vbetool-dpms-off.timer
[Unit]
Description=Timer for VBEtool DPMS off

[Timer]
Unit=betool-dpms-off.service
OnCalendar=ON_CALENDAR

[Install]
WantedBy=timers.target
